# templates/default/tomcat-init.erb
#!/bin/bash
#
# Generated by Chef for <%= node[:fqdn] %>
# Do not edit locally

# Init script for the Apache Tomcat servlet container

# chkconfig: 2345 90 10
# description: Apache Tomcat Server
# processname: <%= @tomcat_service %>
# pidfile: <%= @catalina_pid %>

# Source function library
. /etc/init.d/functions

# Source networking configuration
. /etc/sysconfig/network

# Check that networking is up
[ ${NETWORKING} = "no" ] && exit 0

SERVICE='<%= @tomcat_service.capitalize %>'
LOCKFILE='<%= @tomcat_lock %>'
TOMCAT_USER='<%= @tomcat_user %>'

# Tomcat environment variables
export CATALINA_BASE='<%= @catalina_base %>'
export CATALINA_HOME='<%= @catalina_home %>'
export CATALINA_OPTS='<%= @catalina_opts %>'
export CATALINA_PID='<%= @catalina_pid %>'
export CATALINA_TMPDIR='<%= @catalina_tmp %>'
export JAVA_HOME='<%= @java_home %>'
export JAVA_OPTS='<%= @java_opts %>'

RETVAL=0
start(){
   echo "Starting $SERVICE: "
   su $TOMCAT_USER -c "$CATALINA_HOME/bin/startup.sh"
   RETVAL=$?
   echo
   [ $RETVAL -eq 0 ] && touch $LOCKFILE
   return $RETVAL
}

stop(){
   echo "Shutting down $SERVICE: "
   $CATALINA_HOME/bin/shutdown.sh
   RETVAL=$?
   echo
   [ $RETVAL -eq 0 ] && rm -f $LOCKFILE
   return $RETVAL
}

case "$1" in
  start)
      start
      ;;
   stop)
      stop
      ;;
   restart)
      stop
      start
      ;;
   status)
      status -p <%= @catalina_pid %> tomcat
      ;;
   *)
      echo $"Usage: $0 {start|stop|restart|status}"
      exit 1
      ;;
esac
exit $?
